# packages/p4pi-web/packaging/p4pi-web-genstat.service
[Unit]
Description=Collect statistics for p4edge-web
After=network.target

[Service]
Type=simple
User=root
Group=www-data
WorkingDirectory=/srv/p4edge
ExecStart=/usr/bin/python3 generate-statistics.py

[Install]
WantedBy=multi-user.target
