# packages/p4pi-web/packaging/p4pi-web-dummy-ctrl-plane.service
[Unit]
Description=Collect statistics for p4pi-web
After=network.target

[Service]
Type=simple
User=root
Group=www-data
WorkingDirectory=/srv/p4pi
ExecStart=/usr/bin/python3 dummy_ctrl_plane_connection.py

[Install]
WantedBy=multi-user.target
